# content/etc/init/xbmc-done.conf
start on stopped xbmc-preload or xbmc-failed-start or stopped xbmc-loaded

pre-start exec renice 0 1
