# content/etc/init/xbmc-done.conf
start on stopped xbmc-preload or xbmc-failed-start
