# content/etc/profile.d/xbian-alias.sh
#!/bin/sh

alias xl='tailf /home/xbian/.xbmc/temp/xbmc.log'
alias xll='less /home/xbian/.xbmc/temp/xbmc.log'

